# sstvProxy.service
[Unit]
Description=A simple proxy for SSTV and Plex DVR

[Service]
Environment=
WorkingDirectory=/home/bjzy/sstvProxy/
ExecStart=/home/bjzy/sstvProxy/venv/bin/python /home/bjzy/sstvProxy/sstvProxy.py
Restart=always

[Install]
WantedBy=multi-user.target
